# rootfs/etc/services.d/bom-proxy/run
#!/usr/bin/with-contenv bashio

# Get configuration
export LOG_LEVEL=$(bashio::config 'log_level')
export CACHE_DIR="/data/cache"
export PORT="3000"
export CACHE_TTL_HOURS=$(bashio::config 'cache_ttl_hours')
export TIMESTAMP_REFRESH_INTERVAL=$(bashio::config 'timestamp_refresh_interval')
export MAX_CACHE_SIZE_MB=$(bashio::config 'max_cache_size_mb')

bashio::log.info "Starting BoM Radar Proxy..."
bashio::log.info "Log level: ${LOG_LEVEL}"
bashio::log.info "Cache directory: ${CACHE_DIR}"
bashio::log.info "Cache TTL: ${CACHE_TTL_HOURS} hours"
bashio::log.info "Timestamp refresh interval: ${TIMESTAMP_REFRESH_INTERVAL} seconds"

# Ensure cache directory exists
mkdir -p "${CACHE_DIR}"

# Change to app directory
cd /app || exit 1

# Start the application
exec node server.js
